# roles/influxdb/templates/backup/backup.timer.j2
[Unit]
Description=Run {{ service_name }} backup periodically

[Timer]
OnCalendar={{ backup_cron }}

[Install]
WantedBy=timers.target